# service/emdisplay.service
[Unit]
Description=emdisplay background service

[Service]
ExecStart=/opt/emdisplay/service/service.py

[Install]
WantedBy=multi-user.target